# 2017_hashiconf/present.ini
[program:present]
command=/bin/sh -c "cd /mnt/build && /go/bin/present -http 0.0.0.0:3999 -notes ./hashiconf.slide"
killasgroup=true
stopasgroup=true
